# wallpaper-selector-sync.service
[Unit]
Description=Sync wallpaper from swww to DMS on boot
Before=dms.service graphical-session.target
Wants=swww.service
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/wallpaper-selector sync --verbose
RemainAfterExit=no

[Install]
WantedBy=graphical-session.target
